# supervisord.conf
; 白夜系统 v2.0 — Supervisor 守护配置
; 更新：2026-05-06

[unix_http_server]
file=/tmp/baiye-supervisor.sock
chmod=0700

[supervisord]
logfile=/root/.openclaw/workspace/killer-trading-system/logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=2
loglevel=info
pidfile=/tmp/baiye-supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/baiye-supervisor.sock

; ══════════════════════════════════════════════════════
; Phase5 纸交易引擎 — 每15分钟扫描一次
; ══════════════════════════════════════════════════════
[program:paper_scanner]
command=python3 paper_scanner_loop.py
directory=/root/.openclaw/workspace/killer-trading-system
autostart=true
autorestart=true
startsecs=3
startretries=999
stopwaitsecs=10
killasgroup=true
stdout_logfile=/root/.openclaw/workspace/killer-trading-system/logs/paper_scanner.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/root/.openclaw/workspace/killer-trading-system/logs/paper_scanner_err.log
stderr_logfile_maxbytes=2MB
environment=PYTHONUNBUFFERED="1"
numprocs=1
process_name=%(program_name)s
