# supervisord.conf
; 杀手锏交易系统 — Supervisor 守护配置
; 功能：自动重启、崩溃恢复、日志管理、进程数量限制

[unix_http_server]
file=/tmp/killer-supervisor.sock
chmod=0700

[supervisord]
logfile=/root/.openclaw/workspace/killer-trading-system/logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=2
loglevel=info
pidfile=/tmp/killer-supervisord.pid
nodaemon=false
silent=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/killer-supervisor.sock

; ═══════════════════════════════════════════════════════════
; 纸交易引擎 — 自动重启，最多1个实例
; ═══════════════════════════════════════════════════════════
[program:paper_engine]
command=python3 scripts/paper_engine_v106.py
directory=/root/.openclaw/workspace/killer-trading-system
autostart=true
autorestart=true
startsecs=5
startretries=999
stopwaitsecs=15
killasgroup=true
stopasgroup=true

; 日志
stdout_logfile=/root/.openclaw/workspace/killer-trading-system/logs/paper_engine.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/root/.openclaw/workspace/killer-trading-system/logs/paper_engine_err.log
stderr_logfile_maxbytes=2MB
redirect_stderr=false

; 环境
environment=PYTHONUNBUFFERED="1",PYTHONPATH="/root/.openclaw/workspace/killer-trading-system/scripts"

; 进程数量限制：只允许1个实例（防多开）
numprocs=1
process_name=%(program_name)s
